# bin/rbnb-gumstix
#!/bin/bash
# Init script for the RBNB server on the PacIOOS gumstix deployments
#
# LJM 060607
# The RBNB_HOME location needs to exist and be acceptable to the admin of the
# RBNB host system.
#
# chkconfig: 345 98 98
# description: Start the RBNB server
#

if [ ! "${RBNB_HOME}" ]; then RBNB_HOME=/usr/local/realtime-data; fi
RBNB_LOG='/var/log/rbnb/rbnb.log' # /var/log is volatile on OE/Gumstix
RBNB_PIDFILE='/var/run/rbnb.pid';
RBNB_ARCHIVE_LOC='/var/lib/rbnb';
RBNB_USER='pacioos';
RBNB_SERVER=$(ip addr show wlan0 | grep "inet " | tr -s " " " " | cut -d" " -f3 | cut -d"/" -f1);
RBNB_PORT="3333";
# Amount of memory in MegaBytes for the JVM running RBNB to use
RBNB_MEM='64';
RBNB_NAME='PacIOOSDataTurbine';
REALTIME_DATA='/usr/local/realtime-data';

######################################################################################
start() {
   mkdir -p ${RBNB_ARCHIVE_LOC}
   
   # OE/Gumstix /var/log directory is volatile.  Re-link if necessary
   if [ -d /var/log/rbnb ]; then
     echo "RBNB log directory exists";
   else
     echo "Relinking the log directory.";
     ln -s /usr/local/rbnb/log /var/log/rbnb;
   fi
   
   cd ${RBNB_ARCHIVE_LOC}
   # note rbnb.jar is found in the realtime-data lib directory
   su - ${RBNB_USER}\
    -c "java -Xms${RBNB_MEM}M\
             -Xmx${RBNB_MEM}M\
             -jar ${RBNB_HOME}/lib/rbnb.jar\
             -H ${RBNB_ARCHIVE_LOC}\
             -a ${RBNB_SERVER}:${RBNB_PORT}\
             -n ${RBNB_NAME} > ${RBNB_LOG} 2>&1 &"
   
   RBNB_PID=$(ps -ef | grep rbnb | grep java | tr -s " " " " | cut -d" " -f2);
   echo ${RBNB_PID} > ${RBNB_PIDFILE};

}

######################################################################################
stop() {
  if [ -f ${RBNB_PIDFILE} ]; then
    kill -s SIGTERM `cat ${RBNB_PIDFILE}`;
    rm -f ${RBNB_PIDFILE};
  fi
}
     
case "$1" in
   'start')
       echo "Starting the RBNB server"
       start
       ;;
   'stop')
       echo "Stopping the RBNB server"
       stop
       ;;
   'restart')
       echo "Restarting the RBNB server"
       stop
       start
       ;;
   *)
       echo "Usage: $0 {start|stop|restart} "
       ;;
esac
